# parts/files/leds/boot_40.service
[Unit]
Description=Boot Progress 40%
[Service]
Type=simple
ExecStart=/usr/bin/pinctrl set 5 op dh

[Install]
WantedBy=basic.target
